# files/rc.local
#!/bin/bash
/usr/local/bin/auto_mount_nfs.sh

exit 0